# linux-installation/file-sharer.service
[Unit]
Description=File Sharer - Local Network File Sharing
After=network.target

[Service]
Type=simple
User=%YOUR_USERNAME%
WorkingDirectory=/home/%YOUR_USERNAME%/file-sharer
ExecStart=/usr/bin/python3 /home/%YOUR_USERNAME%/file-sharer/backend/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target